# moj-docker-deploy/apps/templates/upstart_container.conf
description "{{cname}} container"
author "tools@digital.justice.gov.uk"
start on filesystem and started docker
stop on runlevel [!2345]
respawn
respawn limit 5 5
env HOME=/root
script
        {% if cdata.get('name', False) %}
          {% set container_name =  cdata.get('name') %}
        {% else %}
          {% set container_name = cname %}
        {% endif %}

        {% if cdata.get('registry', False) %}
          {% set docker_registry = cdata.get('registry') %}
        {% else %}
          {% set docker_registry = default_registry %}
        {% endif %}


        {% set container_full_name = '%s/%s' % (docker_registry, container_name) %}
        PILLAR_TAG='{{ tag | replace("'", "'\\''") }}'
        if [ -z $TAG ]; then
            TAG=$PILLAR_TAG
        fi

        {%  if 'volumes' in cdata %}
        VOL_OPTS="{% for descr, vol_set in cdata['volumes'].items() %} -v {{vol_set['host']}}:{{vol_set['container']}} {% endfor %}"
        {% endif %}

        {% if 'ports' in cdata %}
        PORT_OPTS="{% for descr, port_set in cdata['ports'].items() %} -p {{port_set['host']}}:{{port_set['container']}} {% endfor %}"
        {% endif %}

        if [ -f /etc/docker_env.d/{{cname}} ]; then
                ENV_OPTS="--env-file /etc/docker_env.d/{{cname}}"
        fi

    docker run --name="{{ cname }}" {{cdata.get('docker_args', '')}} $ENV_OPTS $VOL_OPTS $PORT_OPTS {{container_full_name}}:"$TAG" {{cdata.get('startup_args', '')}}
end script

# Post start script checks that the container it actually running
# It uses a small number of retries since the startup of containers
# is not instantaneous
post-start script
    echo INFO: Checking if container {{cname}} is running
	# Test for success
	for i in {1..5}
	do
		RUNNING_ID="`docker ps -q --filter name={{cname}}`"
		if [ -z ${RUNNING_ID} ]; then
			echo WARNING: Docker not started yet, retrying...
			sleep 3
		else
			echo INFO: Docker instance has started successfully with id ${RUNNING_ID}
			exit 0
		fi
	done
	echo ERROR: Docker failed to start instance
	exit 1
end script

post-stop exec docker rm -f {{ cname }}
